# roles/laptop_power_profile/templates/mem-sleep-default.service.j2
[Unit]
Description=Set mem_sleep default to s2idle when supported
After=multi-user.target
ConditionPathExists=/sys/power/mem_sleep

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'grep -q "s2idle" /sys/power/mem_sleep && echo s2idle > /sys/power/mem_sleep || true'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
