# roles/laptop_power_profile/templates/set-brightness.service.j2
[Unit]
Description=Set laptop backlight brightness to 5 percent
After=display-manager.service graphical.target systemd-user-sessions.service
Wants=display-manager.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'for _ in $(seq 1 15); do brightnessctl -q set 5% || true; sleep 2; done'
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
